# openVas-setup/ssh_services/2tun.service
[Unit]
Description=SSH Tunnel by Adam

[Service]
ExecStart=$scriptspath$/ssh/usedby_openVasgui_tunnel.sh
Restart=always
User=root
RestartSec=3
LimitNOFILE=infinity

[Install]
WantedBy=default.target
