# openVas-setup/ssh_services/2tun.service
[Unit]
Description=SSH Tunnel by Adam

[Service]
ExecStart=/home/kali/adamMADEtunnel/2tunScript.sh
Restart=always
User=kali
RestartSec=3
LimitNOFILE=infinity

[Install]
WantedBy=default.target
